# build_files/waydroid-session.service
[Unit]
Description=SOsharkOS — sessione utente Waydroid (Android) sul display Wayland
Documentation=https://docs.waydro.id/
# La sessione ha bisogno del display Wayland → può partire SOLO dopo il login
# grafico, non al boot puro. Legata alla sessione grafica: se la sessione finisce
# (logout), la sessione Android si ferma con essa.
PartOf=graphical-session.target
After=graphical-session.target
# NB: questa è una USER unit (systemd --user). NON può dipendere da
# waydroid-container.service (SYSTEM unit, altro manager): cross-manager dep non
# esistono. La precondizione "container su" è garantita da androidbox-start, che
# abilita questa unit SOLO dopo che il container è attivo; se al login il
# container è giù, `waydroid session start` esce e il Restart sotto riprova
# blando (no martellamento) finché torna su.

[Service]
Type=simple
# `waydroid session start` resta in foreground finché la sessione Android è viva.
ExecStart=/usr/bin/waydroid session start
ExecStop=/usr/bin/waydroid session stop
# Se il container è giù, riavviare in loop martellerebbe inutilmente: riavvio
# blando con backoff e tetto, mai più di un tentativo ogni 10s.
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=4

[Install]
# NON abilitato di default nell'immagine: lo abilita per-utente `androidbox-start`
# (systemctl --user enable). Quando abilitato, parte al login grafico.
WantedBy=graphical-session.target
